# units/magnets.service
[Unit]
Description=Get magnets

[Service]
Environment=GOPATH=/home/muru/dev
ExecStart=/home/muru/bin/get-magnets
